# fleet/system/commute-evening.service
[Unit]
Description=Commute evening check
After=network.target

[Service]
Type=oneshot
User=root
EnvironmentFile=/opt/lumina-fleet/axon/.env
ExecStart=/usr/bin/python3 /opt/lumina-fleet/commute/commute_check.py evening
SyslogIdentifier=commute-evening
