# fleet/dura/dura-secret-rotation.timer
[Unit]
Description=Dura Secret Rotation Check — weekly on Monday 08:18
Requires=dura-secret-rotation.service

[Timer]
OnCalendar=Mon *-*-* 08:18:00
Persistent=true
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
