# fleet/dura/dura-secret-rotation.service
[Unit]
Description=Dura Secret Rotation Checker
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /opt/lumina-fleet/dura/dura_secret_rotation.py --check
EnvironmentFile=-/opt/lumina-fleet/axon/.env
WorkingDirectory=/opt/lumina-fleet/dura
StandardOutput=journal
StandardError=journal
